# sparse/etc/systemd/system/android-service@vendor.camera-provider-2-4.service.d/10-force-kill.conf
[Service]
Environment=ANDROID_SERVICE_FORCE_KILL=yes
